# home-improvement-planner.service
[Unit]
Description=Home Improvement Planner Web Application
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/home-improvement-planner-nodejs
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=home-improvement-planner
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
